# root/etc/services.d/sockd/run
#!/usr/bin/with-contenv sh

sed -i "s/__IFNAME__/$IFNAME/" /etc/sockd.conf
/usr/sbin/sockd
